# deploy/systemd/news-mcp-dev.service
[Unit]
Description=News MCP Server (dev — mcp_news_dev)
After=network.target

[Service]
User=elucia
WorkingDirectory=/home/elucia/dev/news-mcp
EnvironmentFile=/home/elucia/dev/news-mcp/.env.dev
ExecStart=/home/elucia/dev/news-mcp/.venv/bin/python -m uvicorn news_mcp_server:app --host 127.0.0.1 --port 4821 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
